# udev/35-helper.rules
########################################################################
#
# Description : 35-helper.rules
#
# Authors     : Based on Open Suse Udev Rules
#               kay.sievers@suse.de
#
# Adapted to  : Jim Gifford
# LFS	      : Alexander E. Patrakov
#
# Version     : 00.01
#
# Notes       :
#
########################################################################

# kernel firmware loader
SUBSYSTEM=="firmware",	ACTION=="add", RUN+="firmware.sh"

# check cdrom helper
BUS=="ide",	KERNEL=="hd[a-z]", ACTION=="add", IMPORT="cdrom_id --export $tempnode"
BUS=="scsi",	KERNEL=="sr[0-9]*", ACTION=="add", IMPORT="cdrom_id --export $tempnode"
BUS=="scsi",	KERNEL=="scd[a-z]", ACTION=="add", IMPORT="cdrom_id --export $tempnode"
BUS=="scsi",	KERNEL=="sg[0-9]*", ACTION=="add", DRIVER=="sr", GROUP="cdrom"

BUS=="ide", ENV{ID_CDROM}=="?*",	PROGRAM="cdsymlink_helper.sh %k ide", SYMLINK+="cdrom%c"
BUS=="ide", ENV{ID_CDROM_CD_R}=="?*",	PROGRAM="cdsymlink_helper.sh %k ide", SYMLINK+="cdr%c"
BUS=="ide", ENV{ID_CDROM_CD_RW}=="?*",	PROGRAM="cdsymlink_helper.sh %k ide", SYMLINK+="cdrw%c"
BUS=="ide", ENV{ID_CDROM_DVD}=="?*",	PROGRAM="cdsymlink_helper.sh %k ide", SYMLINK+="dvd%c"
BUS=="ide", ENV{ID_CDROM_DVD_R}=="?*",	PROGRAM="cdsymlink_helper.sh %k ide", SYMLINK+="dvdrw%c"

BUS=="scsi", ENV{ID_CDROM}=="?*",	PROGRAM="cdsymlink_helper.sh %n scsi", SYMLINK+="cdrom%c"
BUS=="scsi", ENV{ID_CDROM_CD_R}=="?*",	PROGRAM="cdsymlink_helper.sh %n scsi", SYMLINK+="cdr%c"
BUS=="scsi", ENV{ID_CDROM_CD_RW}=="?*",	PROGRAM="cdsymlink_helper.sh %n scsi", SYMLINK+="cdrw%c"
BUS=="scsi", ENV{ID_CDROM_DVD}=="?*",	PROGRAM="cdsymlink_helper.sh %n scsi", SYMLINK+="dvd%c"
BUS=="scsi", ENV{ID_CDROM_DVD_R}=="?*",	PROGRAM="cdsymlink_helper.sh %n scsi", SYMLINK+="dvdrw%c"

# adds floppy devices
KERNEL=="fd[0-9]*", ACTION=="add", ATTRS{cmos}=="*", RUN+="create_floppy_devices -c -t $sysfs{cmos} -m %M /dev/%k"
