# udev/35-helper.rules
########################################################################
#
# Description : 35-helper.rules
#
# Authors     : Based on Open Suse Udev Rules
#               kay.sievers@suse.de
#
# Adapted to  : Jim Gifford
# LFS	      : Alexander E. Patrakov
#
# Version     : 00.01
#
# Notes       :
#
########################################################################

# kernel firmware loader
SUBSYSTEM=="firmware",	ACTION=="add", RUN+="firmware_helper"

# check cdrom helper
BUS=="ide",	KERNEL=="hd[a-z]", ACTION=="add", IMPORT="/sbin/cdrom_id --export $tempnode"
BUS=="scsi",	KERNEL=="sr[0-9]*", ACTION=="add", IMPORT="/sbin/cdrom_id --export $tempnode"
BUS=="scsi",	KERNEL="scd[a-z]", ACTION=="add", IMPORT="/sbin/cdrom_id --export $tempnode"
BUS=="scsi",	KERNEL=="sg[0-9]*", ACTION=="add", DRIVER=="sr", GROUP="cdrom"

ENV{ID_CDROM}=="?*",		SYMLINK+="cdrom%m", GROUP="cdrom"
ENV{ID_CDROM_CD_R}=="?*",	SYMLINK+="cdr%m"
ENV{ID_CDROM_CD_RW}=="?*",	SYMLINK+="cdrw%m"
ENV{ID_CDROM_DVD}=="?*",	SYMLINK+="dvd%m"
ENV{ID_CDROM_DVD_R}=="?*",	SYMLINK+="dvdrw%m"

# Run network stuff when the network device is loaded
ACTION=="add", SUBSYSTEM=="net", ENV{INTERFACE}=="*", RUN+="network_helper.sh"
ACTION=="remove", SUBSYSTEM=="net", ENV{INTERFACE}=="*", RUN+="network_helper.sh"

